# opentracker/rootfs/etc/services.d/opentracker/run
#!/usr/bin/with-contenv sh

exec /usr/local/bin/opentracker -f /etc/opentracker/opentracker.conf
